# services/laserbonnet.service
[Unit]
Description=Laserbonnet
After=network-online.target joybonnet.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc /home/pi/src/laserbonnet/scripts/start.sh
WorkingDirectory=/home/pi/src/laserbonnet
Restart=no
User=pi
EnvironmentFile=/home/pi/src/laserbonnet/.env

[Install]
WantedBy=multi-user.target
